# build/systemd/bandwidth-crawler.service
[Unit]
Description="Bandwidth crawler"
After=network-online.target

[Service]
Type=simple
User=crawler
Group=crawler
Restart=always
Environment=PYTHONPATH=./src
WorkingDirectory=/opt/tribler
ExecStart=/usr/bin/python3 scripts/bandwidth_crawler/run_bandwidth_crawler.py --statedir /var/lib/crawler $EXTRA_CRAWLER_ARGS

[Install]
WantedBy=multi-user.target
